# etc/vtllibrary@.service.in
[Unit]
Documentation=man:vtllibrary(1) man:vtlcmd(1)
Description=Robot Library Daemon for Virtual Tape & Robot Library
Before=multi-user.target graphical.target
Requires=mhvtl-load-modules.service
After=mhvtl-load-modules.service
PartOf=mhvtl.service

[Service]
Type=simple
Environment=VERBOSE=0
Environment=DAEMON_DEBUG=""
EnvironmentFile=-@CONF_PATH@/mhvtl.conf
ExecStart=/usr/bin/vtllibrary -F -q%i -v${VERBOSE} ${DAEMON_DEBUG}
ExecStop=/usr/bin/vtlcmd %i exit
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=mhvtl.service
